# src/pawer.service
[Unit]
Description=pawer

[Service]
ExecStart=/home/pithing/Golang/pawer/src/pawer
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
